# files/system_monitor.timer
[Unit]
Description=System monitor timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
